# launchd/com.irchelper.agent-queue.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.irchelper.agent-queue</string>

    <key>ProgramArguments</key>
    <array>
        <string>/Users/kangxi/clawd/agent-queue/agent-queue</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/kangxi/Library/Logs/agent-queue/stdout.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/kangxi/Library/Logs/agent-queue/stderr.log</string>

    <key>WorkingDirectory</key>
    <string>/Users/kangxi/clawd/agent-queue</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>AGENT_QUEUE_PORT</key>
        <string>19827</string>
        <key>AGENT_QUEUE_DISCORD_WEBHOOK_URL</key>
        <string>https://discord.com/api/webhooks/1476188831311986719/YtAzF3zHf7R6Ru6F3HVfrO5hLWVgGSCumZx0pTp0ErEJLFhkGuILAbJOoHW0Bgou7Iwo</string>
        <key>AGENT_QUEUE_DISCORD_USER_ID</key>
        <string>1475329182878928937</string>
    </dict>
</dict>
</plist>
